# chromesearchfix.service.template
[Unit]
Description=ChromeSearchFix

[Service]
ExecStart=__PWD__/run.sh
Restart=always
#Environment=FLASK_APP=search.py
#Environment=DISPLAY=:0

[Install]
WantedBy=basic.target
